# services/supervisord/conf.d/example.conf
; Example Supervisord program configuration
; This file defines a simple Python script that runs under supervisord.

[program:hello_world]
command=python3 /app/scripts/hello_world.py
directory=/app/scripts
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/hello_world.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/hello_world.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=PYTHONUNBUFFERED="1"
user=root
redirect_stderr=true